# agents/aetherbot/com.hyo.aetherbot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.hyo.aetherbot</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/xkai808/Documents/Projects/Hyo/agents/aetherbot/aetherbot.sh</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HYO_ROOT</key>
        <string>/Users/xkai808/Documents/Projects/Hyo</string>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>

    <key>StartInterval</key>
    <integer>900</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/hyo-aetherbot.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/hyo-aetherbot.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/xkai808/Documents/Projects/Hyo</string>
</dict>
</plist>
