# kai/queue/com.hyo.weekly-report.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.hyo.weekly-report</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-c</string>
    <string>HYO_ROOT=/Users/kai/Documents/Projects/Hyo HQ_URL=https://hyo.world HYO_ROOT=/Users/kai/Documents/Projects/Hyo bash /Users/kai/Documents/Projects/Hyo/bin/weekly-report.sh</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key><integer>6</integer>
    <key>Minute</key><integer>0</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>/Users/kai/Documents/Projects/Hyo/kai/ledger/com.hyo.weekly-report.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/kai/Documents/Projects/Hyo/kai/ledger/com.hyo.weekly-report.log</string>
  <key>RunAtLoad</key>
  <false/>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key><string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    <key>HYO_ROOT</key><string>/Users/kai/Documents/Projects/Hyo</string>
    <key>HQ_URL</key><string>https://hyo.world</string>
  </dict>
</dict>
</plist>
