# kai/queue/com.hyo.bore-tunnel.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.hyo.bore-tunnel</string>
  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/bore</string>
    <string>local</string>
    <string>22</string>
    <string>--to</string>
    <string>bore.pub</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/kai/Documents/Projects/Hyo/kai/ledger/bore-tunnel.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/kai/Documents/Projects/Hyo/kai/ledger/bore-tunnel.log</string>
</dict>
</plist>
