# kai/queue/com.hyo.kai-autonomous.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.hyo.kai-autonomous</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-c</string>
    <string>HYO_ROOT=~/Documents/Projects/Hyo HOME=~/  bash ~/Documents/Projects/Hyo/bin/kai-autonomous.sh</string>
  </array>
  <key>StartInterval</key>
  <integer>900</integer>
  <key>RunAtLoad</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/kai/Documents/Projects/Hyo/kai/ledger/kai-autonomous.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/kai/Documents/Projects/Hyo/kai/ledger/kai-autonomous.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>HOME</key>
    <string>/Users/kai</string>
    <key>HYO_ROOT</key>
    <string>/Users/kai/Documents/Projects/Hyo</string>
  </dict>
</dict>
</plist>
