# kai/queue/com.hyo.ticket-sla-enforcer.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.hyo.ticket-sla-enforcer</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-c</string>
    <string>HYO_ROOT=~/Documents/Projects/Hyo bash ~/Documents/Projects/Hyo/bin/ticket-sla-enforcer.sh >> ~/Documents/Projects/Hyo/kai/ledger/ticket-enforcer.log 2>&amp;1</string>
  </array>
  <key>StartInterval</key>
  <integer>1800</integer>
  <key>RunAtLoad</key>
  <false/>
  <key>StandardOutPath</key>
  <string>/Users/kai/Documents/Projects/Hyo/kai/ledger/ticket-enforcer.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/kai/Documents/Projects/Hyo/kai/ledger/ticket-enforcer.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>HOME</key>
    <string>/Users/kai</string>
  </dict>
</dict>
</plist>
