# agents/sam/mcp-server/com.hyo.mcp-tunnel.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.hyo.mcp-tunnel</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>exec $HOME/Documents/Projects/Hyo/agents/sam/mcp-server/tunnel.sh --method cloudflared</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HYO_ROOT</key>
        <string>$HOME/Documents/Projects/Hyo</string>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/hyo-mcp-tunnel.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/hyo-mcp-tunnel.log</string>

    <key>WorkingDirectory</key>
    <string>$HOME/Documents/Projects/Hyo</string>

    <!-- Restart after crash, with exponential backoff -->
    <key>StartInterval</key>
    <integer>10</integer>

    <!-- Allow 10 seconds for startup before retrying -->
    <key>TimeOut</key>
    <integer>10</integer>

    <!-- Limit restart frequency -->
    <key>ThrottleInterval</key>
    <integer>5</integer>
</dict>
</plist>
