# agents/sam/mcp-server/com.hyo.mcp-server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.hyo.mcp-server</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/node</string>
        <string>/Users/xkai808/Documents/Projects/Hyo/agents/sam/mcp-server/server.js</string>
        <string>--transport</string>
        <string>sse</string>
        <string>--port</string>
        <string>3847</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HYO_ROOT</key>
        <string>/Users/xkai808/Documents/Projects/Hyo</string>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/hyo-mcp-server.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/hyo-mcp-server.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/xkai808/Documents/Projects/Hyo</string>
</dict>
</plist>
