# agents/nel/consolidation/com.hyo.simulation.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.hyo.simulation</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>/Users/kai/Documents/Projects/Hyo/bin/dispatch.sh simulate &amp;&amp; /Users/kai/Documents/Projects/Hyo/bin/dispatch.sh simulate-review</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HYO_ROOT</key>
        <string>/Users/kai/Documents/Projects/Hyo</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>23</integer>
        <key>Minute</key>
        <integer>30</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/tmp/hyo-simulation.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/hyo-simulation.log</string>

    <key>WorkingDirectory</key>
    <string>/tmp</string>

    <key>KeepAlive</key>
    <false/>
</dict>
</plist>
