# agents/nel/consolidation/com.hyo.system-maintenance.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.hyo.system-maintenance</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/env</string>
    <string>bash</string>
    <string>-c</string>
    <string>
      export HYO_ROOT="$HOME/Documents/Projects/Hyo";
      export PATH="/usr/local/bin:/opt/homebrew/bin:$PATH";
      if [ -f "$HOME/security/hyo.env" ]; then
        set -a; . "$HOME/security/hyo.env"; set +a;
      fi;
      cd "$HYO_ROOT" || exit 1;
      python3 bin/maintenance-audit.py --all \
        >> "$HYO_ROOT/kai/ledger/maintenance-audit.log" 2>&1
    </string>
  </array>

  <!-- 01:30 MT = 07:30 UTC (MST) or 07:30 UTC (MDT offset handled by system clock) -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>7</integer>
    <key>Minute</key>
    <integer>30</integer>
  </dict>

  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/tmp/com.hyo.system-maintenance.out</string>
  <key>StandardErrorPath</key>
  <string>/tmp/com.hyo.system-maintenance.err</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>/Users/mac</string>
    <key>HYO_ROOT</key>
    <string>/Users/mac/Documents/Projects/Hyo</string>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>
</dict>
</plist>
