# agents/aether/com.hyo.aether.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.hyo.aether</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/kai/Documents/Projects/Hyo/agents/aether/aether.sh</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HYO_ROOT</key>
        <string>/Users/kai/Documents/Projects/Hyo</string>
    </dict>

    <key>StartInterval</key>
    <integer>60</integer>

    <key>StandardOutPath</key>
    <string>/tmp/hyo-aether.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/hyo-aether.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/kai/Documents/Projects/Hyo</string>
</dict>
</plist>
